# 7.3.5-ashamanecore/s6/etc/services.d/worldserver/run
#!/usr/bin/with-contenv bash
cd /ac/bin
exec \
        /usr/bin/screen -D -m -S WORLDSERVER s6-setuidgid abc /ac/bin/worldserver -c /config/${AC_WORLDSERVER_CONFIGFILE}
