# 3.3.5a-mangostwo-playerbots/s6/etc/services.d/mangosd/run
#!/usr/bin/with-contenv bash
cd /mangos/bin
exec \
        /usr/bin/screen -D -m -S MANGOSD s6-setuidgid abc /mangos/bin/mangosd -c /config/${M_MANGOSD_CONFIGFILE}
