# accessory/vpn-ipsec/s6/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bash
UMASK_SET=${UMASK_SET:-022}

umask "$UMASK_SET"

exec \
	s6-setuidgid root /usr/sbin/dnsmasq --no-daemon --no-hosts --no-resolv --conf-file=/etc/dnsmasq.conf
