# 3.3.5a-mangos-two-playerbots/s6/etc/services.d/worldserver/run
#!/usr/bin/with-contenv bash
cd /app/bin
case "${CORE}" in
	"trinitycore")
		exec \
			/usr/bin/screen -D -m -S WORLDSERVER s6-setuidgid abc /app/bin/worldserver -c /config/${WORLDSERVER_CONFIGFILE}
		;;
	"mangos")
		exec \
			/usr/bin/screen -D -m -S WORLDSERVER s6-setuidgid abc /app/bin/mangosd -c /config/${WORLDSERVER_CONFIGFILE}
		;;
	*)
		exit 1
		;;
esac
