# 4.3.4-trinitycore-spp/s6/etc/services.d/worldserver/run
#!/usr/bin/with-contenv bash
cd /tc/bin
exec \
        /usr/bin/screen -D -m -S WORLDSERVER s6-setuidgid abc /tc/bin/worldserver -c /config/${TC_WORLDSERVER_CONFIGFILE}
